# nopanel/docker/laravel/docker-compose/supervisord/conf.d/supervisord.conf
[supervisord]
nodaemon=true
loglevel=warn
logfile_maxbytes=0
logfile=/dev/stderr

[program:mail]
command=php /root/natsclient/src/email.php
autostart=true
autorestart=true
startretries=999999999999999999999
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:mail2]
command=php /root/natsclient/src/email2.php
autostart=true
autorestart=true
startretries=999999999999999999999
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:mail3]
command=php /root/natsclient/src/email3.php
autostart=true
autorestart=true
startretries=999999999999999999999
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:mail4]
command=php /root/natsclient/src/email4.php
autostart=true
autorestart=true
startretries=999999999999999999999
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:mail5]
command=php /root/natsclient/src/email5.php
autostart=true
autorestart=true
startretries=999999999999999999999
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0


[program:mail6]
command=php /root/natsclient/src/email6.php
autostart=true
autorestart=true
startretries=999999999999999999999
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0


[program:mail7]
command=php /root/natsclient/src/email7.php
autostart=true
autorestart=true
startretries=999999999999999999999
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0


[program:mail8]
command=php /root/natsclient/src/email8.php
autostart=true
autorestart=true
startretries=999999999999999999999
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0


[program:mail9]
command=php /root/natsclient/src/email9.php
autostart=true
autorestart=true
startretries=999999999999999999999
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0



[program:mail10]
command=php /root/natsclient/src/email10.php
autostart=true
autorestart=true
startretries=999999999999999999999
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:member]
command=php /root/natsclient/src/member.php
autostart=true
autorestart=true
startretries=999999999999999999999
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:member2]
command=php /root/natsclient/src/member2.php
autostart=true
autorestart=true
startretries=999999999999999999999
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:member3]
command=php /root/natsclient/src/member3.php
autostart=true
autorestart=true
startretries=999999999999999999999
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:member4]
command=php /root/natsclient/src/member4.php
autostart=true
autorestart=true
startretries=999999999999999999999
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:member5]
command=php /root/natsclient/src/member5.php
autostart=true
autorestart=true
startretries=999999999999999999999
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:member6]
command=php /root/natsclient/src/member6.php
autostart=true
autorestart=true
startretries=999999999999999999999
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:member7]
command=php /root/natsclient/src/member7.php
autostart=true
autorestart=true
startretries=999999999999999999999
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:member8]
command=php /root/natsclient/src/member8.php
autostart=true
autorestart=true
startretries=999999999999999999999
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:member9]
command=php /root/natsclient/src/member9.php
autostart=true
autorestart=true
startretries=999999999999999999999
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:member10]
command=php /root/natsclient/src/member10.php
autostart=true
autorestart=true
startretries=999999999999999999999
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0


[program:transaksi]
command=php /root/natsclient/src/transaksi.php
autostart=true
autorestart=true
startretries=999999999999999999999
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:transaksi2]
command=php /root/natsclient/src/transaksi2.php
autostart=true
autorestart=true
startretries=999999999999999999999
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:transaksi3]
command=php /root/natsclient/src/transaksi3.php
autostart=true
autorestart=true
startretries=999999999999999999999
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:transaksi4]
command=php /root/natsclient/src/transaksi4.php
autostart=true
autorestart=true
startretries=999999999999999999999
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:transaksi5]
command=php /root/natsclient/src/transaksi5.php
autostart=true
autorestart=true
startretries=999999999999999999999
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:transaksi6]
command=php /root/natsclient/src/transaksi6.php
autostart=true
autorestart=true
startretries=999999999999999999999
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:transaksi7]
command=php /root/natsclient/src/transaksi7.php
autostart=true
autorestart=true
startretries=999999999999999999999
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:transaksi8]
command=php /root/natsclient/src/transaksi8.php
autostart=true
autorestart=true
startretries=999999999999999999999
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:transaksi9]
command=php /root/natsclient/src/transaksi9.php
autostart=true
autorestart=true
startretries=999999999999999999999
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:transaksi10]
command=php /root/natsclient/src/transaksi10.php
autostart=true
autorestart=true
startretries=999999999999999999999
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0


[program:wa]
command=php /root/natsclient/src/wa.php
autostart=true
autorestart=true
startretries=999999999999999999999
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0


[program:wa2]
command=php /root/natsclient/src/wa2.php
autostart=true
autorestart=true
startretries=999999999999999999999
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:wa3]
command=php /root/natsclient/src/wa3.php
autostart=true
autorestart=true
startretries=999999999999999999999
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0



